# infra/systemd/learncoding-offsite-sync.service
[Unit]
Description=Codestead verified encrypted offsite publication
Documentation=/opt/learncoding/docs/runbooks/backup-and-restore.md
After=network-online.target learncoding-backup.service
Wants=network-online.target
RequiresMountsFor=/mnt/learncoding-backups
OnFailure=learncoding-alert@%n.service

[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/bash /opt/learncoding/scripts/backup/offsite-sync.sh
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=full
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
UMask=0077
TimeoutStartSec=4h
TimeoutStopSec=2m
