# infra/systemd/learncoding-compose.service
[Unit]
Description=LearnCoding trusted application stack
Documentation=file:/opt/learncoding/docs/deployment.md
After=docker.service network-online.target local-fs.target libvirtd.service learncoding-runner-firewall.service
Requires=docker.service learncoding-runner-firewall.service
Wants=network-online.target libvirtd.service
RequiresMountsFor=/opt/learncoding /etc/learncoding /srv/learncoding
OnFailure=learncoding-alert@%n.service
StartLimitIntervalSec=15min
StartLimitBurst=5

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/learncoding
EnvironmentFile=/etc/learncoding/compose.env
ExecStart=/usr/bin/env PATH=/usr/sbin:/usr/bin:/sbin:/bin /usr/bin/bash /opt/learncoding/infra/ops/start-production-stack.sh --startup-wait 600
ExecReload=/usr/bin/env PATH=/usr/sbin:/usr/bin:/sbin:/bin /usr/bin/bash /opt/learncoding/infra/ops/start-production-stack.sh --startup-wait 600
ExecStop=/usr/bin/env -i HOME=/nonexistent LANG=C LC_ALL=C PATH=/usr/sbin:/usr/bin:/sbin:/bin DOCKER_CONFIG=/nonexistent DOCKER_HOST=unix:///var/run/docker.sock COMPOSE_PROJECT_NAME=learncoding COMPOSE_PROFILES= /usr/bin/docker compose --project-name learncoding --env-file /etc/learncoding/compose.env -f /opt/learncoding/compose.yaml down --remove-orphans
Restart=on-failure
RestartSec=15s
TimeoutStartSec=15min
TimeoutStopSec=5min

[Install]
WantedBy=multi-user.target
