# infra/systemd/learncoding-restore-drill.service
[Unit]
Description=LearnCoding supervised isolated restore drill
Documentation=/opt/learncoding/docs/runbooks/backup-and-restore.md
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target
RequiresMountsFor=/mnt/learncoding-backups
OnFailure=learncoding-alert@%n.service

# This unit is intentionally manual: attach the offline age identity, set
# AGE_IDENTITY_FILE in /etc/learncoding/backup.env, then start it explicitly.
[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/bash /opt/learncoding/scripts/backup/restore-drill.sh
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=read-only
ProtectSystem=full
UMask=0077
TimeoutStartSec=4h
