# infra/systemd/learncoding-backup.timer
[Unit]
Description=Run LearnCoding encrypted backup nightly

[Timer]
OnCalendar=*-*-* 02:15:00
RandomizedDelaySec=20m
Persistent=true
AccuracySec=1m
Unit=learncoding-backup.service

[Install]
WantedBy=timers.target
