# infra/systemd/learncoding-runner-guest-firewall.service
[Unit]
Description=Codestead runner guest ingress firewall
Documentation=file:/opt/learncoding/docs/runbooks/runner-isolation.md
After=local-fs.target
Before=network-pre.target learncoding-runner.service
Wants=network-pre.target

[Service]
Type=oneshot
User=root
Group=root
UMask=0077
ExecStartPre=/usr/sbin/nft --check --file /opt/learncoding/infra/runner-vm/guest-runner.nft
ExecStart=/usr/sbin/nft --file /opt/learncoding/infra/runner-vm/guest-runner.nft
ExecReload=/usr/sbin/nft --file /opt/learncoding/infra/runner-vm/guest-runner.nft
RemainAfterExit=yes
NoNewPrivileges=true
CapabilityBoundingSet=CAP_NET_ADMIN
AmbientCapabilities=CAP_NET_ADMIN
PrivateDevices=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
RestrictAddressFamilies=AF_NETLINK AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
