# infra/systemd/learncoding-production-load-gate.service
[Unit]
Description=Codestead supervised production load gate
Documentation=file:/opt/learncoding/docs/runbooks/load-testing.md
Requires=learncoding-production-load-control.service
After=learncoding-production-load-control.service network-online.target
Wants=network-online.target
RequiresMountsFor=/opt/learncoding /etc/learncoding /var/lib/learncoding-production-load-evidence
ConditionFileIsExecutable=/usr/bin/node
AssertPathExists=/opt/learncoding/node_modules/tsx/package.json
OnFailure=learncoding-alert@%n.service

[Service]
Type=exec
User=learncoding-load-gate
Group=learncoding-load-gate
WorkingDirectory=/opt/learncoding
EnvironmentFile=/etc/learncoding/production-load.env
Environment=LOAD_MODE=production
Environment=LOAD_ALLOW_REMOTE=1
Environment=LOAD_SCOPE=codestead-project-only
Environment=LOAD_PROJECT=learncoding
Environment=LOAD_DISPOSABLE_FAULTS_CONFIRMED=1
Environment=LOAD_EVIDENCE_ROOT=/var/lib/learncoding-production-load-evidence
Environment=LOAD_ACTIVE_RELEASE_PATH=/etc/learncoding/active-release.env
Environment=LOAD_CONTROL_SOCKET=/run/learncoding/load-control.sock
ExecStartPre=/usr/bin/bash /opt/learncoding/infra/ops/validate-production-load-host-runtime.sh
ExecStart=/usr/bin/node --import tsx /opt/learncoding/scripts/load-smoke.ts
Restart=no
RuntimeMaxSec=5h
TimeoutStartSec=2min
TimeoutStopSec=2min
KillSignal=SIGTERM
FinalKillSignal=SIGKILL
SendSIGKILL=yes
UMask=0077
LimitCORE=0
CapabilityBoundingSet=
AmbientCapabilities=
NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectHostname=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
LockPersonality=yes
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/var/lib/learncoding-production-load-evidence
InaccessiblePaths=-/etc/learncoding/secrets -/var/lib/learncoding-production-load -/run/docker.sock -/run/libvirt -/dev/kvm
